# deploy/systemd/camlab.service.d/shutdown-splash.conf
# SPDX-FileCopyrightText: Copyright (C) 2026, UAB Kurokesu
# SPDX-License-Identifier: GPL-3.0-or-later
#
# Arm the shutdown splash while Cage still runs, so no tty text flashes.
# ExecStop runs before systemd kills the main process.

[Service]
ExecStop=+/usr/local/bin/camlab-shutdown-splash
