# deploy/camlab.service
# SPDX-FileCopyrightText: Copyright (C) 2026, UAB Kurokesu
# SPDX-License-Identifier: GPL-3.0-or-later
#
# Rendered by scripts/setup/service.sh (CAMLAB_USER / CAMLAB_REPO_DIR are
# substituted). Runs the camlab GUI fullscreen under a Cage kiosk on tty1.

[Unit]
Description=camlab camera bench tool (Kurokesu)
After=systemd-user-sessions.service getty@tty1.service plymouth-quit-wait.service
Conflicts=getty@tty1.service

[Service]
Type=simple
User=CAMLAB_USER
PAMName=camlab
WorkingDirectory=CAMLAB_REPO_DIR

# Writable, persistent state owned by the service user (creates /var/lib/camlab
# and exports $STATE_DIRECTORY). Holds the boot-persistent mode/fps selection.
# Phase 5 (read-only root) must keep this path on a writable data partition,
# excluded from the overlay, so the selection survives reboots.
StateDirectory=camlab
StateDirectoryMode=0750

TTYPath=/dev/tty1
StandardInput=tty-fail
StandardOutput=journal
StandardError=journal
SyslogIdentifier=camlab

Environment=XCURSOR_SIZE=1
# Deliberately do NOT raise LIBCAMERA_LOG_LEVELS: ERROR lines (e.g. AR0822
# embedded-data parse failures) are the integrity signal we surface.
Environment=PYTHONUNBUFFERED=1

ExecStart=/usr/bin/cage -- /usr/bin/python3 -m camlab

Restart=on-failure
RestartSec=2

[Install]
WantedBy=multi-user.target
